# systemd/nova-orchestrator.service
[Unit]
Description=Nova Orchestrator - Cascade Task Execution Engine
After=network.target
Documentation=https://github.com/yourusername/nova

[Service]
Type=simple
User=%i
# Update WorkingDirectory to point to your project directory
# All paths (CASCADE.md, cascade_state.json, logs/) will be relative to this directory
WorkingDirectory=%h/nova
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 -m program_nova.engine.orchestrator --daemon
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/nova/orchestrator.log
StandardError=append:/var/log/nova/orchestrator.error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/nova

[Install]
WantedBy=multi-user.target
